# Entorno/ProvisionamientoEVO/modules/tomcat/templates/tomcat.erb
#!/bin/bash
#
# tomcat
#
# chkconfig: 345 99 01
# description:  Start up the Tomcat servlet engine.

# Source function library.
. /etc/init.d/functions


RETVAL=$?
CATALINA_HOME="/u01/app/apache/apache-tomcat-<%= @tomcat_version %>"

case "$1" in
 start)
        if [ -f $CATALINA_HOME/bin/startup.sh ];
          then
        echo $"Starting Tomcat"
#            /bin/su tomcat $CATALINA_HOME/bin/startup.sh
            sudo $CATALINA_HOME/bin/startup.sh			
        fi
    ;;
 stop)
        if [ -f $CATALINA_HOME/bin/shutdown.sh ];
          then
        echo $"Stopping Tomcat"
#            /bin/su tomcat $CATALINA_HOME/bin/shutdown.sh
            sudo $CATALINA_HOME/bin/shutdown.sh
        fi
    ;;
 *)
    echo $"Usage: $0 {start|stop}"
    exit 1
    ;;
esac

exit $RETVAL